# .bw/pacmir/arch-linux/mirror/usr/lib/systemd/system/pacmir-mirror-rsync.service
[Unit]
Description=rsync arch linux mirrors by time

[Service]
Type=oneshot
ExecStart=/usr/bin/true